# standalone/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

; ── Elasticsearch ─────────────────────────────────────────────────────────
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
user=elasticsearch
environment=ES_JAVA_OPTS="%(ENV_ES_JAVA_OPTS)s"
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

; ── Wait for ES before starting proxy and MCP ─────────────────────────────
[program:es-wait]
command=bash -c "echo '[es-wait] Waiting for Elasticsearch...' && until curl -sf http://127.0.0.1:9200/_cluster/health >/dev/null 2>&1; do sleep 2; done && echo '[es-wait] Elasticsearch is ready' && supervisorctl start proxy mcp"
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20

; ── mitmproxy ─────────────────────────────────────────────────────────────
[program:proxy]
command=mitmdump -s proxy_addon.py --set listen_port=8080
directory=/app
autostart=false
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=30

; ── MCP server ────────────────────────────────────────────────────────────
[program:mcp]
command=python mcp_server.py --transport streamable-http
directory=/app
autostart=false
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=30
